# overlay/etc/init.d/00_link_mount
#!/sbin/openrc-run

command="/usr/bin/link_mount start || true"
