# overlay/etc/init.d/usb_gadget
#!/sbin/openrc-run

command="/usr/bin/usb_gadget start || true"
